# com.cloudflare.tailor-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cloudflare.tailor-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>tailor-api</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/madhukonda</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/cloudflared-tailor-api.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cloudflared-tailor-api.error.log</string>
</dict>
</plist>

